# templates/base.timer.template
[Unit]
Description={{DESCRIPTION}}

[Timer]
OnCalendar={{SCHEDULE}}
Persistent=true

[Install]
WantedBy=timers.target
